# conf/etc/systemd/system/card_codex.socket
[Unit]
Description=card codex gunicorn socket

[Socket]
ListenStream=/run/card_codex/socket

[Install]
WantedBy=sockets.target
